# launchd/com.r3lay.services.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  r3LAY LaunchAgent — edit paths before installing.
  Install: bash launchd/install.sh
-->
<dict>
    <key>Label</key>
    <string>com.r3lay.services</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- EDIT THIS PATH to match your r3LAY install location -->
        <string>__R3LAY_DIR__/launchd/r3lay-up.sh</string>
        <string>--detach</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/r3lay-launcher.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/r3lay-launcher.err</string>
</dict>
</plist>
